# deploy/systemd/tokio.service
# tokio.service — web Next.js em 127.0.0.1:3002 (VPS compartilhada, ADR 0007)
# Instalação (operador, com sudo): copiar para /etc/systemd/system/tokio.service
# Node via nvm: o ExecStart usa bash -l para carregar o nvm do usuário tokio.

[Unit]
Description=Tokio web (Next.js)
After=network.target

[Service]
User=tokio
Group=tokio
WorkingDirectory=/home/tokio/Tokio/web
EnvironmentFile=/home/tokio/Tokio/.env
Environment=PORT=3002
Environment=HOSTNAME=127.0.0.1
Environment=NODE_ENV=production
# standalone build: o deploy copia .next/static e public para dentro de .next/standalone
ExecStart=/usr/bin/env bash -lc 'exec node .next/standalone/server.js'
Restart=always
RestartSec=3
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
