# deploy/systemd/tokio-engine.service
# tokio-engine.service — supervisor do engine (VPS compartilhada, ADR 0007)
# Mantém gateway e runners como processos separados com restart
# individual. Tudo binda em 127.0.0.1 (deploy/engine-processes.yaml).
# Instalação (operador, com sudo): copiar para /etc/systemd/system/ e
# ADICIONAR ao sudoers do usuário tokio (restart/status), junto do tokio.service.

[Unit]
Description=Tokio engine (gateway + runners)
After=network.target

[Service]
User=tokio
Group=tokio
WorkingDirectory=/home/tokio/Tokio
EnvironmentFile=/home/tokio/Tokio/.env
ExecStart=/home/tokio/Tokio/.venv/bin/python -m engine.supervisor
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=25
MemoryMax=1536M
CPUQuota=150%

[Install]
WantedBy=multi-user.target
